# puppet/modules/daemontools/files/daemontools.conf
# daemontools

start on stopped rc

stop on shutdown

respawn
exec /command/svscanboot
